# battery-tray.service
[Unit]
Description=Battery Tray UPS Monitor
After=graphical.target network.target
Wants=graphical.target

[Service]
Type=simple
User=root
WorkingDirectory=INSTALL_DIR
ExecStart=INSTALL_DIR/.venv/bin/python3 INSTALL_DIR/batteryTray.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=/var/run/lightdm/root/:0

[Install]
WantedBy=graphical.target
