# systemd/battery-monitor.service
[Unit]
Description=Battery Monitor — Waveshare UPS HAT D
Documentation=https://github.com/waveshare/UPS-HAT-D
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
# Adjust paths to match where you installed the application
WorkingDirectory=%h/battery-monitor
ExecStart=/usr/bin/python3 %h/battery-monitor/app.py --config %h/battery-monitor/config/config.yaml
Restart=on-failure
RestartSec=10
# Prevent restart storms — give up after 5 rapid failures in 60 seconds
StartLimitBurst=5
StartLimitIntervalSec=60

# Allow access to the I2C device
SupplementaryGroups=i2c

# Resource limits
# Nice value: +5 yields to interactive processes while maintaining responsiveness
Nice=5
# OOM killer weight: slightly prefer to kill this over user processes
OOMScoreAdjust=100

# Environment
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
# Wayland support (if GNOME uses Wayland)
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Logging — journald captures stderr automatically
StandardOutput=journal
StandardError=journal
SyslogIdentifier=battery-monitor

[Install]
WantedBy=graphical-session.target
